# server/com.webcam-magic.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- LaunchAgent: the CORS/Host-rewrite proxy (server/proxy.js) on :11435 in
     front of Ollama. ngrok points at 11435. setup.sh substitutes __PROXY_JS__. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.webcam-magic.proxy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>__PROXY_JS__</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/webcam-magic-proxy.log</string>
  <key>StandardErrorPath</key><string>/tmp/webcam-magic-proxy.log</string>
</dict>
</plist>
